# homepage/.config/containers/systemd/homepage.container
[Unit]
Description=Homepage dashboard
After=network-online.target

[Container]
Image=ghcr.io/gethomepage/homepage:latest
ContainerName=homepage
Environment=PUID=1000 PGID=1000
Environment=HOMEPAGE_ALLOWED_HOSTS=home.lab,homepage.lab,homepage.lab.albandiguer.dev
PublishPort=3030:3000

# Persistent volume on mounted hcloud volume
Volume=${QUADLET_STORAGE_PATH}/homepage:/app/config

# Rootless Podman socket for container integrations (read-only)
Volume=%t/podman/podman.sock:/var/run/docker.sock:ro

# Health check
HealthCmd=nc -z localhost 3000 || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=30s

[Service]
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target default.target
