# plane/.config/containers/systemd/plane-proxy.container
[Unit]
Description=Plane Reverse Proxy
After=network-online.target plane-web.service plane-api.service plane-space.service plane-admin.service plane-live.service
Requires=plane.pod plane-web.service plane-api.service
After=plane-pod.service
BindsTo=plane-pod.service

[Container]
Image=docker.io/makeplane/plane-proxy:stable
ContainerName=plane-proxy
Pod=plane.pod

# Proxy configuration
Environment=FILE_SIZE_LIMIT=5242880
Environment=BUCKET_NAME=plane
Environment=LISTEN_HTTP_PORT=80
Environment=LISTEN_HTTPS_PORT=443
Environment=SITE_ADDRESS=:80

# Persist nginx config and SSL data
Volume=${QUADLET_STORAGE_PATH}/plane/proxy-config:/config
Volume=${QUADLET_STORAGE_PATH}/plane/proxy-data:/data

[Service]
Restart=always
RestartSec=30
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target
