# pihole/.config/containers/systemd/pihole.container
[Unit]
Description=Pi-hole DNS Server
After=network-online.target

[Container]
Image=docker.io/pihole/pihole:latest
ContainerName=pihole
# For rootless, we'll use host network mode since we need port 53
# This works because we set net.ipv4.ip_unprivileged_port_start=53
Network=host
Secret=pihole_webpassword,type=env,target=FTLCONF_webserver_api_password

# Timezone
Environment=TZ=UTC

# Pi-hole FTL configuration
Environment=FTLCONF_dns_listeningMode=all
Environment=FTLCONF_dns_upstreams=1.1.1.1;1.0.0.1
Environment=FTLCONF_misc_dnsmasq_lines=address=/.lab/10.0.0.1;address=/.lab.albandiguer.dev/10.0.0.1;address=/.home.arpa/10.0.0.1
Environment=FTLCONF_webserver_port=8053

# Volumes for persistent data on configured storage volume
Volume=${QUADLET_STORAGE_PATH}/pihole/etc:/etc/pihole
Volume=${QUADLET_STORAGE_PATH}/pihole/dnsmasq.d:/etc/dnsmasq.d

# Health check
HealthCmd=/usr/bin/dig +short +norecurse +retry=0 @127.0.0.1 pi.hole
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target default.target